# dist/tootctl-media-remove.timer
[Unit]
Description=Runs tootctl-media-remove.service

[Timer]
OnCalendar=*-*-* 03:00:00

[Install]
WantedBy=multi-user.target